# quadlet/postgres.container
[Unit]
Description=TimescaleDB (PostgreSQL) for Feelinq

[Container]
Image=docker.io/timescale/timescaledb:latest-pg17
Environment=POSTGRES_USER=feelinq
Environment=POSTGRES_PASSWORD=feelinq
Environment=POSTGRES_DB=feelinq
Volume=%h/feelinq-data/postgres:/var/lib/postgresql/data:Z
Network=feelinq.network

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
